# conf/64-btrfs.rules
RUN+="/bin/mknod /dev/btrfs-control c 10 234"
