# conf/rockstor-hdparm.service
Description=Rockstor hdparm settings
After=suspend.target

[Service]
Type=oneshot
ExecStart=

[Install]
WantedBy=suspend.target basic.target