# conf/rockstor.service
[Unit]
Description=Rockstor startup script
After=rockstor-pre.service
Requires=rockstor-pre.service

[Service]
Environment="DJANGO_SETTINGS_MODULE=settings"
WorkingDirectory=/opt/rockstor
ExecStart=/root/.local/bin/poetry run supervisord -c /opt/rockstor/etc/supervisord.conf

[Install]
WantedBy=multi-user.target
