# conf/rockstor.service.in
[Unit]
Description=RockStor startup script
After=rockstor-pre.service
Requires=rockstor-pre.service

[Service]
ExecStart=${buildout:depdir}/bin/supervisord -c ${buildout:depdir}/etc/supervisord.conf

[Install]
WantedBy=multi-user.target
